# service/metro-server5002/run
#!/bin/sh
exec 2>&1

PERL=/home/metro/perl5/perlbrew/perls/perl-Metro-5.16.3/bin/perl
DAEMON=/home/metro/Metro/script/metro-server5002

exec \
  envdir ./env \
  setuidgid metro \
  $PERL $DAEMON
